# stegovault_1.0.0/etc/systemd/system/stegovault.service
[Unit]
Description=StegoVault Digital Steganography Vault
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/stegovault
Environment=PORT=9001
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/stegovault/.venv/bin/python app.py
Restart=always

[Install]
WantedBy=multi-user.target
