# stegovault_1.0.0/DEBIAN/postinst
#!/bin/bash
set -e

APP_DIR="/opt/stegovault"
VENV_DIR="$APP_DIR/.venv"

# Ensure the app directory exists (should be handled by dpkg)
if [ ! -d "$APP_DIR" ]; then
    echo "Error: Application directory $APP_DIR not found."
    exit 1
fi

echo "Setting up StegoVault virtual environment (dynamic)..."
python3 -m venv "$VENV_DIR"
"$VENV_DIR/bin/pip" install --upgrade pip
"$VENV_DIR/bin/pip" install -r "$APP_DIR/requirements.txt"

echo "Configuring permissions..."
chown -R root:root "$APP_DIR"
chmod -R 755 "$APP_DIR"

echo "Reloading systemd and starting stegovault.service..."
systemctl daemon-reload
systemctl enable stegovault.service
systemctl restart stegovault.service

echo ""
echo "----------------------------------------------------------------"
echo "StegoVault successfully installed and running on port 9001!"
echo "Service name: stegovault.service"
echo "Launcher: /usr/local/bin/stegovault-launcher"
echo "----------------------------------------------------------------"
echo ""

